# src/includes.chroot/home/user/.config/autostart/gnome-setup.desktop
[Desktop Entry]
Name=Initial gnome setup for Open Semantic Desktop Search
Type=Application
NoDisplay=true
Exec=/home/user/gnome-setup.sh
